# luci-app-netmonitor/root/etc/uci-defaults/luci-netmonitor
#!/bin/sh

sed -i 's|^;DatabaseDir.*|DatabaseDir "/etc/vnstat"|' /etc/vnstat.conf

bash /www/netmonitor/vnstati/vnstati.sh >/dev/null 2>&1

exit 0
